# mixpost_horizon.conf
[program:mixpost_horizon]
process_name=%(program_name)s
command=php /var/www/artisan horizon  # Adjust path if needed
autostart=true
autorestart=true
user=developer  # Set to the user running Laravel in Docker
numprocs=1
startsecs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/horizon.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stopwaitsecs=5
stopsignal=SIGTERM
stopasgroup=true
killasgroup=true
